# packaging/debian/postinst
#!/bin/bash
# Post-installation script for npu-proxy

set -e

case "$1" in
    configure)
        PACKAGE_VERSION="$(dpkg-query -W -f='${Version}' npu-proxy 2>/dev/null | sed 's/-[^-]*$//')"
        WHEELHOUSE="/usr/share/npu-proxy/wheels"
        CONSTRAINTS_FILE="/usr/share/npu-proxy/constraints.txt"
        APP_WHEEL="$(find "$WHEELHOUSE" -maxdepth 1 -name "npu_proxy-${PACKAGE_VERSION}-*.whl" | head -n 1)"

        # Create npu-proxy user if it doesn't exist
        if ! getent passwd npu-proxy > /dev/null 2>&1; then
            adduser --system --group --no-create-home --home /var/lib/npu-proxy npu-proxy
        fi
        
        # Create model cache directory
        mkdir -p /var/lib/npu-proxy/models
        chown -R npu-proxy:npu-proxy /var/lib/npu-proxy
        
        # Create log directory
        mkdir -p /var/log/npu-proxy
        chown -R npu-proxy:npu-proxy /var/log/npu-proxy
        
        # Install the packaged runtime into a dedicated virtual environment
        mkdir -p /opt/npu-proxy
        python3 -m venv /opt/npu-proxy/venv
        /opt/npu-proxy/venv/bin/python -m ensurepip --upgrade
        if [ -z "$APP_WHEEL" ]; then
            echo "Unable to find the packaged npu-proxy wheel for version $PACKAGE_VERSION" >&2
            exit 1
        fi
        /opt/npu-proxy/venv/bin/python -m pip install --constraint "$CONSTRAINTS_FILE" --force-reinstall --upgrade "$APP_WHEEL"
        
        # Reload systemd
        systemctl daemon-reload
        
        echo ""
        echo "NPU Proxy installed successfully!"
        echo ""
        echo "To start the service:"
        echo "  sudo systemctl start npu-proxy"
        echo ""
        echo "To enable at boot:"
        echo "  sudo systemctl enable npu-proxy"
        echo ""
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
